# .config/systemd/user/syncthing.service
[Unit]
Description=Syncthing
ConditionPathIsMountPoint=/mnt/media

[Service]
ExecStart=/usr/bin/syncthing

[Install]
WantedBy=default.target
